# conf/htpc/etc/pm/sleep.d/10-persistent-ir.rules
KERNEL=="event*", ATTRS{idProduct}=="0613", ATTR{idVendor}="0471", SYMLINK="input/mceremote", RUN+="/usr/local/sbin/ir-keytable.sh"

